# scripts/systemd/nocturne.service
[Unit]
Description=Nocturne - Autonomous Coding Orchestrator (daemon)
After=network-online.target
Wants=network-online.target
StartLimitBurst=5
StartLimitIntervalSec=300

[Service]
Type=simple
ExecStart=%h/projects/nocturne/.venv/bin/nocturne daemon
WorkingDirectory=%h/projects/nocturne
Restart=on-failure
RestartSec=30s
EnvironmentFile=-%h/.config/nocturne/env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=nocturne

[Install]
WantedBy=default.target
